# 企业微信/py.cron
0 */2 * * * cd /src && python3 check_cookie.py
30 9 * * * cd /src && python3 post_zichan.py
