# 企业微信/py.cron
9      7       *       *       *       cd /src && python send_msg.py
